# projectroot/var/spool/cron/crontabs/root
* * * * *  /usr/sbin/logrotate -s /var/log/logrotate.status /etc/logrotate.conf
